# mongobak.cron
0 0 * * * root /usr/local/ornu/mongo-backup/mongobak.py
